# init/lxc-ci.reboot.service
[Unit]
Description=lxc-ci - Reboot

[Service]
Type=oneshot
WorkingDirectory=/lxc-ci
ExecStart=/lxc-ci/bin/reboot
