# init/lxc-ci.kernel-check.service
[Unit]
Description=lxc-ci - kernel check

[Service]
Type=oneshot
WorkingDirectory=/lxc-ci
ExecStart=/bin/sh -c "exec /lxc-ci/bin/kernel-check"
