# init/lxc-ci.update-cache.service
[Unit]
Description=lxc-ci - Update cache
Requires=network-online.target
After=network-online.target
Before=lxc-ci.ready.service

[Service]
Type=oneshot
WorkingDirectory=/lxc-ci
ExecStart=/bin/sh -c "exec /lxc-ci/bin/update-cache"

[Install]
WantedBy=multi-user.target
