# init/lxc-ci.cleanup.service
[Unit]
Description=lxc-ci - Cleanup

[Service]
Type=oneshot
WorkingDirectory=/lxc-ci
ExecStart=/lxc-ci/bin/cleanup
