# init/lxc-ci.update-cache.timer
[Unit]
Description=lxc-ci - Update cache timers

[Timer]
OnBootSec=6h
OnUnitActiveSec=6h

[Install]
WantedBy=timers.target
