# init/lxc-ci.update-cache.timer
[Unit]
Description=lxc-ci - Update cache timers

[Timer]
OnBootSec=1h
OnUnitActiveSec=1h

[Install]
WantedBy=timers.target
