# init/lxc-ci.cleanup.timer
[Unit]
Description=lxc-ci - Cleanup timer

[Timer]
OnBootSec=1h
OnUnitActiveSec=1h

[Install]
WantedBy=timers.target
