# init/lxc-ci.cleanup.service
[Unit]
Description=lxc-ci - Cleanup

[Service]
Type=oneshot
WorkingDirectory=/lxc-ci
ExecStart=/bin/sh -c "exec /lxc-ci/bin/cleanup"
