# init/lxc-ci.reboot.timer
[Unit]
Description=lxc-ci - Reboot timer

[Timer]
OnBootSec=1d

[Install]
WantedBy=timers.target
