# init/lxc-ci.oom-check.timer
[Unit]
Description=lxc-ci - OOM check timer

[Timer]
OnBootSec=5m
OnUnitActiveSec=5m

[Install]
WantedBy=timers.target
