# init/lxc-ci.git-refresh.timer
[Unit]
Description=lxc-ci - Git repository refresh timer

[Timer]
OnBootSec=1h
OnUnitActiveSec=1h

[Install]
WantedBy=timers.target
